# deploy/ohlc-refresher.service
[Unit]
Description=Edge Catcher — Coinbase OHLC Refresher
Documentation=https://github.com/goldenwo/edge-catcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=edge-catcher
Group=edge-catcher
WorkingDirectory=/opt/edge-catcher
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/edge-catcher/.env
ExecStart=/opt/edge-catcher/.venv/bin/python -m edge_catcher ohlc-refresh --config config.local/paper-trader-spotfair.yaml
StandardOutput=append:/var/log/edge-catcher/ohlc-refresher.log
StandardError=append:/var/log/edge-catcher/ohlc-refresher.log

# Restart always — the refresher is a long-running poll loop; idempotent upsert
# makes a hard kill safe (no graceful drain needed).
Restart=always
RestartSec=5
StartLimitIntervalSec=600
StartLimitBurst=10

# Light daemon: one HTTP poll + small upsert every ~20s.
MemoryMax=128M
CPUQuota=20%
Nice=10

[Install]
WantedBy=multi-user.target
