# deploy/paper-trader-15m.service
[Unit]
Description=Edge Catcher Paper Trader 15m
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=private-infra
WorkingDirectory=/home/private-infra/edge-catcher
EnvironmentFile=/home/private-infra/.private-infra/.env
ExecStart=/usr/bin/python3 -m edge_catcher paper-trade-15m --db data/paper_trades.db
Restart=always
RestartSec=30
Nice=10
StandardOutput=append:/var/log/edge-catcher/paper-trader-15m.log
StandardError=append:/var/log/edge-catcher/paper-trader-15m.log
MemoryMax=256M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
